# overlay/system/etc/init/setup-swap.conf
description "Setup 4GB Swap File"

# We wait for android partitions to be mounted so we have access to /userdata
start on mounted-android

task

script
    SWAPFILE="/userdata/swapfile"
    
    # Check if swap is already active
    if grep -q "$SWAPFILE" /proc/swaps; then
        exit 0
    fi

    # Create swap file if it doesn't exist
    if [ ! -f "$SWAPFILE" ]; then
        # Try fallocate first (fast), fallback to dd
        fallocate -l 4G "$SWAPFILE" || dd if=/dev/zero of="$SWAPFILE" bs=1M count=4096
        chmod 600 "$SWAPFILE"
        mkswap "$SWAPFILE"
    fi
    
    # Enable swap
    swapon "$SWAPFILE"
end script
